# cron/purge.cron
0   0   *   *   *       find /var/www/htdocs/feuille -type f -mtime +7 -exec rm {} +
