# cron/startup.cron
@reboot                 feuille -U https://my.paste.bin
